# ocr-service.service
[Unit]
Description=PaddleOCR Service for expenseApp
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ocr-service
ExecStart=/usr/bin/python3 /opt/ocr-service/ocr_service.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment="PYTHONUNBUFFERED=1"
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

[Install]
WantedBy=multi-user.target

