# deployment/backup/expenseapp-backup.service
[Unit]
Description=Expense App PostgreSQL Backup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/expenseapp/backup.env
ExecStart=/usr/bin/bash /opt/expenseapp/deployment/backup/backup.sh
User=root
Group=root


