# deployment/backup/expenseapp-backup.timer
[Unit]
Description=Run Expense App Backup Daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target


