# deployment/backend/expenseapp-backend.service
[Unit]
Description=Expense App Backend API
After=network.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/expenseapp/backend.env
WorkingDirectory=/opt/expenseapp/backend
ExecStart=/usr/bin/node /opt/expenseapp/backend/dist/server.js
Restart=always
RestartSec=5
User=expenseapp
Group=expenseapp
AmbientCapabilities=
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/expenseapp/uploads
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target


